# site-cookbooks/cloudenabling/files/default/uwsgi.conf
# Generated by Chef, do not modify
# simple uWSGI script
# for /etc/init/uwsgi.conf
# COMMANDS:
# start uwsgi
# stop uwsgi

description "uwsgi tiny instance"
start on runlevel [2345]
stop on runlevel [06]

exec /opt/cloudenabling/current/bin/uwsgi parts/uwsgi/uwsgi.xml --master --processes 4 --die-on-term --socket :49152 --module cloudenabling_wsgi --chdir /opt/cloudenabling/current --ignore-sigpipe --limit-post 1024000000 --post-buffering 4096 --ignore-sigpipe --limit-post 1024000000 >> /var/log/cloudenabling/uwsgi-1.log 2>&1
